# com.ralph.tracker.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.ralph.tracker</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Library/Frameworks/Python.framework/Versions/3.14/bin/python3</string>
        <string>/Users/mahesh/Documents/plg-stock-analysis-1/ralph_tracker.py</string>
        <string>--passive</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/mahesh/Documents/plg-stock-analysis-1</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>

    <key>ThrottleInterval</key>
    <integer>60</integer>

    <key>StandardOutPath</key>
    <string>/Users/mahesh/Documents/plg-stock-analysis-1/ralph_launchd.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/mahesh/Documents/plg-stock-analysis-1/ralph_launchd_error.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/Library/Frameworks/Python.framework/Versions/3.14/bin</string>
    </dict>
</dict>
</plist>
